# packaging/runtime/neuwerk.service.in
[Unit]
Description=Neuwerk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=__RUNTIME_BINARY_DIR__/neuwerk-bootstrap
ExecStart=__RUNTIME_BINARY_DIR__/neuwerk-launch
Restart=on-failure
RestartSec=5
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_SYS_ADMIN CAP_IPC_LOCK
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_SYS_ADMIN CAP_IPC_LOCK
LimitMEMLOCK=infinity
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/neuwerk /var/lib/neuwerk /var/log /run /tmp __RUNTIME_PREFIX__

[Install]
WantedBy=multi-user.target
